# files/system/etc/profile.d/lowlatency.sh
#!/bin/bash
export AMD_DEBUG=lowlatencyenc
